# files/etc/init.d/easydrcom
#!/bin/sh /etc/rc.common
# easydrcom script

START=99

start(){
	easydrcom -b -r -c /etc/config/easydrcom.conf &
}

stop(){
	killall easydrcom
}
